# ubuntu_brightbox_ruby/files/etc/profile.d/Z99-ruby-gems-env.sh
#!/bin/bash

export GEM_PATH="$HOME/gems"
export GEM_HOME="$HOME/gems"
export PATH="$PATH:$HOME/gems/bin"
